# fritz-repeater-start-radio.cron
*/5 * * * * sh /root/fritz-repeater-start-radio.sh '<PASSWORD>' 2
