# ssh-tunnel.cron
*/5 * * * * sh ssh-tunnel.sh <destination> [port]
