# ssh-tunnel.cron
*/5 * * * * sh ssh-tunnel <destination> [port]
