# ssh-tunnel.cron
*/1 * * * * sh ssh-tunnel.sh <destination> [port]
